# services/ir-event-listening.service
[Unit]
Description=Consume infrared cell events
After=ir-setup.service

[Service]
ExecStart=/usr/local/bin/argosir

[Install]
WantedBy=default.target
