# services/setup-ir.service
[Unit]
Description=Configure the infrared cell

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-ir.sh

[Install]
WantedBy=default.target
